# config/includes.chroot_after_packages/etc/systemd/user/start-installer.service
[Service]
ExecStart=/bin/sh -c 'pkexec calamares'
